# templates/systemd.service.j2
# This file was generated by Ansible for {{ansible_fqdn}}
# Do NOT modify this file by hand!

{% set app = item -%}
{% set work_dir = app.work_dir|default(celery_work_dir) -%}
{% set log_dir = app.log_dir|default(celery_log_dir) -%}
{% set item_queue = app.queue|default('') -%}
{% set item_name = "%s%s%s" % (app.action|default('worker'), item_queue, '-%s' % item_queue if item_queue else '') -%}
{% set item_user = app.user|default(celery_user) -%}
{% set item_group = app.group|default(celery_group) -%}
{% set run_dir = app.run_dir|default(celery_run_dir) -%}
{% set celery_app = app.app_module|default(celery_app_module) -%}
{% set celery_log_file = '%s/%s.log' % (log_dir, item_name) -%}
{% set celery_log_level = app.log_level|default(celery_log_level) -%}
{% set celery_pid_file = '%s/%s.pid' % (run_dir, item_name) -%}
{% set celery_queue_arg = '-q %s' % item_queue if item_queue else '' -%}
{% set celery_concurrency_arg = '-c %s' % app.concurrency|default(celery_concurrency) if app.action == 'worker' else '' -%}
{% set env = app.env|default(celery_env) -%}
{% set env_file = app.env_file|default(celery_env_file) -%}

[Unit]
Description="Celery workers for {{ celery_app_name }}"
After=syslog.target

[Service]
{% if item_user == 'root' %}
Environment='C_FORCE_ROOT=true' 'HOME=/root'
{% endif %}

{% if env -%}
Environment={% for key in env %}'{{key}}={{env[key]}}' {% endfor %}
{% endif -%}
{% if env_file -%}
EnvironmentFile={{ env_file }}
{% endif -%}

WorkingDirectory={{ work_dir }}
User=root

{# Run any pre-exec commands we have defined, e.g `. /path/to/your/env_vars/.env`#}
{% for line in app.celery_pre_exec|default(celery_pre_exec) %}
ExecStartPre={{ line }}
{% endfor %}

ExecStart={{ celery_bin }} {{ app.action }} {{ celery_queue_arg }} {{ celery_concurrency_arg }} --pidfile="{{ celery_pid_file }}" --logfile="{{ celery_log_file }}" --loglevel="{{ celery_log_level }}" --app="{{ celery_app }}" --workdir="{{ work_dir }}" --uid="{{ item_user }}" --gid="{{ item_group }}" {% if app.action == 'beat' %}--schedule={{ celery_beat_dbfile }}{% endif %}

Restart=always
KillSignal=SIGQUIT

[Install]
WantedBy=multi-user.target
